# linux_p2/crontab
*/5 * * * * q /home/q/git/devops_edu/linux_p2/script.sh

